# build/files/etc/init.d/v2ray_redirect
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

USE_PROCD=1
START=99

service_triggers() {
  procd_add_reload_trigger v2ray
}

start_service() {
  config_load v2ray
  config_get_bool enabled "enabled" "enabled" "0"
  uci set firewall.v2ray.enabled="$enabled"
  uci commit firewall
  /etc/init.d/firewall reload
}

stop_service() {
  config_load v2ray
  uci set firewall.v2ray.enabled="0"
  uci commit firewall
  /etc/init.d/firewall reload
}
