# provision/tmp/provision/uci-defaults/01-set-password
#
# Check if this is not the very first boot but a firmware upgrade
[[ $(uci -q get system.@system[0].firstboot_done) ]] && exit 0

PASSWORD=$(uci -q get provision.@provision[0].password)
[ ! "$PASSWORD" ] && PASSWORD="hoogmoon"

# Set password for root account

echo -e "${PASSWORD}\n${PASSWORD}" | passwd

# Activate TTY login

uci set system.@system[0].ttylogin="1"
uci commit system

exit 0
