# build/files/etc/uci-defaults/92-firewall-pptp
#
# Check if this is not the very first boot but a firmware upgrade
[[ $(uci -q get system.@system[0].firstboot_done) ]] && exit 0

# Add TCP 1723, and GRE to INPUT ALLOW

rule_name=$(uci add firewall rule)
uci batch << EOI
set firewall.$rule_name.enabled='1'
set firewall.$rule_name.target='ACCEPT'
set firewall.$rule_name.src='wan'
set firewall.$rule_name.src_ip='192.168.0.0/16'
set firewall.$rule_name.proto='tcp'
set firewall.$rule_name.dest_port='1723'
set firewall.$rule_name.name='Allow-PPTP-Input'
EOI

rule_name=$(uci add firewall rule)
uci batch << EOI
set firewall.$rule_name.enabled='1'
set firewall.$rule_name.target='ACCEPT'
set firewall.$rule_name.src='wan'
set firewall.$rule_name.src_ip='192.168.0.0/16'
set firewall.$rule_name.proto='gre'
set firewall.$rule_name.name='Allow-GRE-Input'
EOI

uci commit firewall

exit 0
