# build/files/etc/uci-defaults/92-lan-ip-table-lookup-100
#
# Check if this is not the very first boot but a firmware upgrade
[[ $(uci -q get system.@system[0].firstboot_done) ]] && exit 0

# Configure LAN/PPP traffic to lookup
# ip route table 100 first
# Used for OpenVPN gateway redirection
# without altering main routing table

# Add 'LAN' zone rule

rule_name=$(uci add network rule)
uci batch << EOI
set network.$rule_name.lookup='100'
set network.$rule_name.in='lan'
set network.$rule_name.priority='31000'
EOI

# Add PPP rule for every interface

for intf in ppp0 ppp1 ppp2 ppp3; do
rule_name=$(uci add network rule)
uci batch << EOI
set network.$rule_name.lookup='100'
set network.$rule_name.in=$intf
set network.$rule_name.priority='31001'
EOI
done

uci commit network

exit 0
