# build/files/etc/profile.d/96-bash.sh
[ "$(command -v shopt)" ] && shopt -s histappend
HISTFILE=/tmp/.bash_history
HISTCONTROL=ignoreboth 
HISTFILESIZE=10000
HISTSIZE=1000

export LESSHISTFILE=-
